# git-config-article/etc/cron.hourly/syncer
#!/bin/bash

bash /root/configurations/syncer.sh /root/configurations /etc /usr/local/etc

